# systemd/.config/systemd/user/backup.service
[Unit]
Description=Backup Config 

[Service]
ExecStart=/home/sven/.local/scripts/backup.sh

[Install]
WantedBy=default.target
