# systemd/.config/systemd/user/backup.timer
[Unit]
Description=Backup configs every two hours.

[Timer]
OnCalendar=00/1:00:00

[Install]
WantedBy=timers.target
