# systemd/.config/systemd/user/battery.service
[Unit]
Description=Battery Status

[Service]
ExecStart=/home/sven/.local/scripts/battery.sh

[Install]
WantedBy=default.target
