# systemd/.config/systemd/user/battery.timer
[Unit]
Description=Check battery every minute.

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
